# piframe-qt/services/piframe-sync.timer
[Unit]
Description=PiFrame Photo Sync Timer (hourly)
Requires=piframe-sync.service

[Timer]
# Run on boot after 2 minutes
OnBootSec=2min
# Run every hour
OnUnitActiveSec=1h
# Randomize by up to 5 minutes to avoid network congestion
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
